# config/upstart.conf
description "Elva IRC bot"

start on starting ruby
stop on stopping ruby
respawn

exec sudo -E -u deploy sh -c "cd /var/www/elva/current && bundle exec foreman run ./bin/elva"
